# src/main/resources/external/user-profile-api.service
[Unit]
Description=user-profile-api: User Profile API

[Service]
Type=forking

# WARNING: environment.rc and config files can override any variable defined here.
# WARNING: Do not guess, but read wrapper.sh script instead.

Environment=JAVA_HOME="/usr/lib/jvm/java-17-amazon-corretto"
Environment=APP_JAVA_OPTS="-server -Xms500m -Xmx1000m -XX:+HeapDumpOnOutOfMemoryError -Xlog:gc*:file=%h/user-profile-api/log/user-profile-api.GC.log"
Environment=CLASSPATH=""
Environment=APP_ENTRY_POINT="-jar %JAR_FILE%"
Environment=APP_ARGS=""
Environment=SPRING_CONFIG_ADDITIONAL_LOCATION="file:%h/user-profile-api/config.local.yaml"
Environment=APP_HOME="%h/user-profile-api/current"
Environment=APP_BASE="%h/user-profile-api"

WorkingDirectory=%h/user-profile-api/current
ExecStart=/bin/bash -c "\"%h/bin/user-profile-api\" start"
ExecStop=/bin/bash -c "\"%h/bin/user-profile-api\" stop"

[Install]
WantedBy=default.target
